# 自动化情报收集系统/templates/com.local.ai-intel-radar.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.ai-intel-radar.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/tanglengyi/Documents/自动化情报收集系统/scripts/run_daily.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/tanglengyi/Documents/自动化情报收集系统</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/tanglengyi/Documents/自动化情报收集系统/data/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tanglengyi/Documents/自动化情报收集系统/data/launchd.err.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
